# deploy/systemd/ultrabot-weekly.service
[Unit]
Description=Ultra Weekly Digest
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ultrabot
EnvironmentFile=/opt/ultrabot/.env
ExecStart=/opt/ultrabot/.venv/bin/python /opt/ultrabot/services/weekly_digest.py
Restart=always
RestartSec=30
User=ultrabot

[Install]
WantedBy=multi-user.target